# etc/skel/.config/autostart/arcobobo-welcome.desktop
[Desktop Entry]
Type=Application
Exec=/usr/bin/arcobobo-welcome
Hidden=false
X-GNOME-Autostart-enabled=true
Name=ArcoBobo Welcome
Comment=
